# scheduler/com.hermes.morning-feishu-pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.hermes.morning-feishu-pipeline</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/gino/project_ai_trading/scheduler/morning_feishu_pipeline.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/gino/project_ai_trading</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>50</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/gino/project_ai_trading/logs/morning_feishu_pipeline.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/gino/project_ai_trading/logs/morning_feishu_pipeline.err.log</string>

  <key>RunAtLoad</key>
  <false/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONPATH</key>
    <string>/Users/gino/project_ai_trading</string>
  </dict>
</dict>
</plist>
